# optional/ytdl-archiver.service
[Unit]
Description=ytdl-archiver YouTube Archiving Service
After=network.target

[Service]
Type=simple

# ---------------------------------------------------------------------------
# REQUIRED: Update these values for your system
# ---------------------------------------------------------------------------

# Path to your ytdl-archiver project directory
WorkingDirectory=/home/username/ytdl-archiver

# Your system user (run 'whoami' to find your username)
User=your_username
Group=your_username

# ---------------------------------------------------------------------------
# OPTIONAL: Customize runtime arguments
# ---------------------------------------------------------------------------

# Example: Use custom playlists file
# ExecStart=ytdl-archiver archive -p /home/username/ytdl-archiver/playlists.toml

# Example: Use custom archive directory
# ExecStart=ytdl-archiver archive -d /path/to/your/video/archive

# Example: Use custom config file
# ExecStart=ytdl-archiver -c /path/to/config.toml archive

# Default: Uses config at ~/.config/ytdl-archiver/config.toml and playlists at ~/.config/ytdl-archiver/playlists.toml
ExecStart=ytdl-archiver archive

# Restart automatically if it crashes
Restart=on-failure

[Install]
WantedBy=multi-user.target
